# src/python/joypad/bin/service/blued/run
#!/bin/bash
exec /usr/local/libexec/bluetooth/bluetoothd
