# install/raspbian/rootfs/opt/cannybots/runtime/services/nodered/run
#!/bin/sh
exec /opt/cannybots/runtime/nodered/cannybox/bin/www

