# install/raspbian/rootfs/opt/cannybots/runtime/services/nodered/run
#!/bin/sh
cd /opt/cannybots/runtime/nodered/cannybox/
sleep 2
exec /usr/local/bin/node --max-old-space-size=128 ./bin/www

