# install/raspbian/rootfs/opt/cannybots/runtime/services/camera/run
#!/bin/sh
cd /opt/cannybots/runtime/camera/
exec ./start.sh


