# install/raspbian/rootfs/opt/cannybots/runtime/services/bluetoothd/run
#!/bin/sh
sleep2
exec /lib/bluetooth/bluetoothd


