# install/raspbian/rootfs/opt/cannybots/runtime/services/lescan/run
#!/bin/sh
sleep 5

exec hcitool lescan --discovery=l

