# install/raspbian/rootfs/opt/cannybots/runtime/services/scratch/run
#!/bin/sh
cd /opt/cannybots/runtime/scratch/
sleep 2
exec ./start.sh


