# install/raspbian/rootfs/opt/cannybots/runtime/services/scratch/run
#!/bin/sh
export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/games:/usr/games
sleep 30
exec python -m cannybots.scratchagent


